# aiem-user.service
[Unit]
Description=aiem Web UI (skills & MCP manager)
After=default.target

[Service]
Type=simple
ExecStart=%h/.local/bin/aiem serve --host 127.0.0.1 --port 8787
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
